# autostart/signage-view.desktop
[Desktop Entry]
Type=Application
Name=Signage Viewer (feh)
X-GNOME-Autostart-enabled=true
Exec=/bin/sh -lc 'sleep 3; unclutter -idle 1 -root & feh --quiet --fullscreen --hide-pointer --reload 5 --auto-zoom "$HOME/signage/frame.png"'
